# led-clock.cron
* * * * * root /usr/local/bin/led-clock
